# image.d/etc/crontabs/seeddms
*/10 * * * * /srv/seeddms/seeddms/utils/seeddms-indexer >/dev/null
