# image.d/etc/services.d/php-fpm/run
#!/bin/execlineb -P

/bin/export LD_PRELOAD /usr/lib/preloadable_libiconv.so
/usr/sbin/php-fpm7 -F
