# image.d/etc/services.d/cron/run
#!/bin/execlineb -P

/bin/fdmove -c 2 1
/usr/sbin/crond -f -d 8
